# launchd/com.user.autocrm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.autocrm</string>

  <key>RunAtLoad</key>
  <true/>

  <!-- Run periodically while user is logged in and Mac is awake -->
  <key>StartInterval</key>
  <integer>120</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>AUTOCRM_NOTION_TOKEN</key>
    <string>REPLACE_WITH_NOTION_INTEGRATION_TOKEN</string>
    <key>AUTOCRM_NOTION_DATABASE_ID</key>
    <string>REPLACE_WITH_PEOPLE_DATABASE_ID</string>

    <!-- Optional: only needed if you use Beeper Desktop -->
    <key>BEEPER_DESKTOP_TOKEN</key>
    <string>REPLACE_WITH_BEEPER_TOKEN</string>

    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>-m</string>
    <string>autocrm.main</string>
  </array>

  <key>StandardOutPath</key>
  <string>/tmp/autocrm.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/autocrm.err</string>
</dict>
</plist>
